# rootfs/etc/services.d/php5-fpm/run
#!/usr/bin/env sh

exec php5-fpm --nodaemonize -e