# rootfs/etc/services.d/logs/run
#!/usr/bin/env sh

exec tail -qF /data/logs/*